# digital-signage.service
[Unit]
Description=Digital Signage Display Server
After=network.target

[Service]
Type=notify
User=pi
Group=pi
WorkingDirectory=/opt/digital-signage
Environment="PATH=/opt/digital-signage/venv/bin"
ExecStart=/opt/digital-signage/venv/bin/gunicorn --config /opt/digital-signage/gunicorn_config.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
