# digital-signage.service
[Unit]
Description=Digital Signage Display Server
After=network.target

[Service]
Type=notify
User=pi
Group=pi
WorkingDirectory=/opt/digital-signage
Environment="PATH=/opt/digital-signage/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/digital-signage/venv/bin/gunicorn --config /opt/digital-signage/gunicorn_config.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStartSec=30
TimeoutStopSec=10
PrivateTmp=true
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
